# unraid/config/plugins/user.scripts/customSchedule.cron
# Generated cron schedule for user.scripts
5,10,15,20,25,30,35,40,45,50,55 * * * * /usr/local/emhttp/plugins/user.scripts/startCustom.php /boot/config/plugins/user.scripts/scripts/vm-powerprofile/script > /dev/null 2>&1

